# deploy/aicrm-inbox-queue-runtime.service
[Unit]
Description=AI-CRM Webhook Inbox PostgreSQL Queue Runtime (claimless standby)
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=AICRM_QUEUE_WORKER_GENERATION=0
ExecStart=/bin/bash -lc 'export AICRM_LISTENER_DATABASE_URL="${AICRM_LISTENER_DATABASE_URL:-${DATABASE_URL:-}}"; test -n "$AICRM_LISTENER_DATABASE_URL"; source /home/ubuntu/venvs/openclaw/bin/activate; exec python scripts/run_execution_runtime.py --queue-kind webhook --generation 0'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
