# deploy/aicrm-queue-soak-snapshot.service
[Unit]
Description=AI-CRM Queue Runtime Soak Evidence Snapshot
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate; exec python scripts/ops/manage_queue_runtime_soak.py --action snapshot --actor systemd-soak-snapshot --reason periodic-15-minute-evidence'
