# deploy/aicrm-queue-invariant-check.timer
[Unit]
Description=Run AI-CRM queue invariant checks every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=15
Unit=aicrm-queue-invariant-check.service

[Install]
WantedBy=timers.target
