# deploy/openclaw-wecom-callback-inbox-worker.service
[Unit]
Description=OpenClaw WeCom Callback Inbox Worker - process queued callback ingress
After=network.target openclaw-wecom-postgres.service
Requires=openclaw-wecom-postgres.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/极简 crm
EnvironmentFile=/home/ubuntu/.openclaw-wecom-pg.env
Environment=AICRM_WECOM_CALLBACK_INBOX_WORKER_EXECUTE=1
Environment=AICRM_WECOM_CALLBACK_INBOX_WORKER_BATCH_SIZE=20
Environment=AICRM_WECOM_CALLBACK_INBOX_WORKER_MAX_EXECUTE_BATCH_SIZE=20
ExecStart=/bin/bash -lc 'source /home/ubuntu/venvs/openclaw/bin/activate && python scripts/run_wecom_callback_inbox_worker.py --execute --limit ${AICRM_WECOM_CALLBACK_INBOX_WORKER_BATCH_SIZE:-20}'
