# deploy/openclaw-ai-audience-scheduler.timer
[Unit]
Description=Write AI Audience daily refresh intents at 02:00 Asia/Shanghai

[Timer]
OnCalendar=*-*-* 02:00:00 Asia/Shanghai
Persistent=true
Unit=openclaw-ai-audience-scheduler.service

[Install]
WantedBy=timers.target
